# targeting.service
[Unit]
Desctiption=Targeting Bootloader

[Service]
ExecStart=/home/option16pi/NetworkTableClient/startup.sh

[Install]
WantedBy=multi-user.target
